# supervisor.conf
[program:blog_app]
command=/var/www/django_blog/venv/bin/gunicorn django_blog:wsgi -c gunicorn.conf.py
directory=/var/www/django_blog
autostart=true
autorestart=true